# provisioning/roles/7dtd_common/files/etc/systemd/system/7dtd_sync_start.service
[Unit]
Description = 7dtd save data backup to s3

[Service]
Type = simple
ExecStart = /root/.bin/save_data_sync.sh start

[Install]
WantedBy = multi-user.target
